# request-catcher.conf
#!upstart

description "request catcher server"
author "Jackson Owens"

start on (local-filesystems and net-device-up IFACE=eth0)
stop on shutdown

respawn
respawn limit 5 60

script
  cd $GOPATH/src/github.com/jbowens/requestcatcher
  go build
  exec sudo -u requestcatcher ./requestcatcher 104.131.146.52 8000
end script
